# vendor/plugins/mysql_s3_backup_tool/sample.cron
#!/bin/bash
ruby /etc/mysql_backup/mysql_tools.rb backup > /dev/null
exit 0
